# rootfs/etc/s6-overlay/s6-rc.d/init-presence_tracker/up
/etc/s6-overlay/s6-rc.d/init-presence_tracker/run
